# phase2/configs/backend.service
[Unit]
Description=Notes Backend API Service
After=syslog.target network.target postgresql.service

[Service]
Type=simple
User=${USER}
WorkingDirectory=${PROJECT_ROOT}/phase1/app/backend

# Java command to run the JAR file
ExecStart=/usr/bin/java -jar ${PROJECT_ROOT}/phase1/app/backend/target/notes-backend-1.0.0.jar

# Auto restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=notes-backend

[Install]
WantedBy=multi-user.target
